# work7_Systemd/w2/httpd.target
[Unit]
Wants=httpd@one.service httpd@two.service
