# work7_Systemd/w1/log_otus.timer
[Unit]
Description=timet log Kostyuk_Ruslan

[Timer]
OnCalendar=*:*:0,30

#OnBootSec=30sec
#OnUnitActiveSec=1d


[Install]
WantedBy=timers.target
